# openclaw-manager.service
[Unit]
Description=OpenClaw Manager Service
After=network.target

[Service]
Type=simple
ExecStart=%h/.openclaw/projects/manager/src/bin/managerd --config %h/.openclaw-manager/config.toml --static-dir %h/.openclaw/projects/manager/src/frontend/dist
Restart=on-failure
RestartSec=5
Environment=OPENCLAW_JWT_SECRET=
WorkingDirectory=%h/.openclaw/projects/manager/src

[Install]
WantedBy=default.target
