# service/openclaw-manager.service
[Unit]
Description=OpenClaw Manager Service
After=network.target

[Service]
Type=simple
Environment=PATH=${node_bin_dir}:${openclaw_bin_dir}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=%h/.openclaw-manager/bin/managerd --config %h/.openclaw-manager/config/config.toml --static-dir %h/.openclaw-manager/web
Restart=on-failure
RestartSec=5
Environment=OPENCLAW_JWT_SECRET=
WorkingDirectory=%h/.openclaw-manager

[Install]
WantedBy=default.target
